# linux/snapgear-2.6-p42/examples/fb_cons/etc/inittab
::sysinit:/etc/init.d/rcS

ttyS0::respawn:/bin/sh

::respawn:/bin/sh
