# linux/snapgear-2.6-p42/romfs/etc/inittab
# Basic Init Table
#

# System initialization (runs when system boots).
::sysinit:/etc/init.d/rcS

::respawn:/bin/sh
